# sys-apps/logwatch/files/logwatch.cron
#!/bin/sh

if [ -x /usr/sbin/logwatch ] ; then
	exec /usr/sbin/logwatch --output mail
fi
